# softwareupdate_service/K96Rpi_update.timer
[Unit]
Description=K96Rpi Software Update Service Timer

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target